# supervisor_configs/engageai_core.conf
[program:engageai_core]
command=/home/bo/projects/engageAI_v2/.venv/bin/gunicorn engageai_core.wsgi -c /home/bo/projects/engageAI_v2/supervisor_configs/gunicorn.conf.py
directory=/home/bo/projects/engageAI_v2/
user=bo

autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=10
stopasgroup=true
killasgroup=true

# log with supervisor
stdout_logfile=/home/bo/projects/engageAI_v2/logs/supervisor_stdout.log
stderr_logfile=/home/bo/projects/engageAI_v2/logs/supervisor_stderr.log
redirect_stderr=false
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

environment=PYTHONPATH="/home/bo/projects/engageAI_v2/engageai_core:/home/bo/projects/engageAI_v2/",DJANGO_SETTINGS_MODULE="engageai_core.settings",DJANGO_ENV="production"
# environment=DJANGO_SETTINGS_MODULE=knowledge_base.settings
